# packaging/launchd/com.local.livis-relayd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- 加载前将 __BUN__、__PROJECT_DIR__、__HOME__ 替换为绝对路径。 -->
  <key>Label</key>
  <string>com.local.livis-relayd</string>
  <key>ProgramArguments</key>
  <array>
    <string>__BUN__</string>
    <string>run</string>
    <string>src/index.ts</string>
    <string>serve</string>
    <string>--config</string>
    <string>__HOME__/.livis-relay/config.json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>PATH</key>
    <string>__HOME__/.local/bin:__HOME__/.bun/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <!-- plist 整数 63 等于八进制 077。 -->
  <key>Umask</key>
  <integer>63</integer>
  <key>StandardOutPath</key>
  <string>__HOME__/.livis-relay/daemon.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.livis-relay/daemon.stderr.log</string>
</dict>
</plist>
